# backend/scripts/com.wolfcapital.kite-token-sync.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Copy both plists + kite_token_sync_job.sh to your machine, fix REPO paths
     in the shell script, then:
       launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.wolfcapital.kite-token-sync.plist
       launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.wolfcapital.kite-token-catchup.plist
     See backend/README.md § Kite token sync. -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.wolfcapital.kite-token-sync</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/YOU/stock_ai/backend/scripts/kite_token_sync_job.sh</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>5</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>5</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>5</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>5</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>5</integer></dict>
  </array>
</dict>
</plist>
